# dev/docker/magento/crontab
* * * * * php /var/www/magento/cron.php
